# configuration/etc/cron.d/mediadownloader
#!/bin/bash
*/15 * * * * root cd /home/djohnson/apps/mediadownloader && /home/djohnson/apps/java/bin/java -jar mediadownloader.jar /etc/mediaserver.properties
